# turk/skel/config/drivers.conf
[program:debug-driver]
command = drivers/debug.py
process_name=%(program_name)s
autostart=true
redirect_stderr=true

